# example.service
[Unit]
Description=run test example

[Service]
Type=simple
User=root
WorkingDirectory=/tmp/
ExecStart=example
PIDFile=/var/run/example.pid
Restart=always

[Install]
WantedBy=multi-user.target
